# kafka/supervisor/metricbeat.conf
[program:metricbeat]
command=service metricbeat start
autostart=true
autorestart=true
